# vllm.service
# vllm.service -- Systemd user service for vLLM + DiffusionGemma 26B A4B
# DGX Spark (GB10, Grace Blackwell)
#
# Prerequisites (headless/SSH-only):
#   sudo apt install dbus dbus-user-session
#   sudo loginctl enable-linger $USER
#   echo 'export DBUS_SESSION_BUS_ADDRESS=\"unix:path=/run/user/$(id -u)/bus\"' >> ~/.bashrc
#
# Installation (as normal user -- no root):
#   mkdir -p ~/.config/systemd/user/
#   cp vllm.service ~/.config/systemd/user/
#   chmod +x /srv/vllm2/vllm-server.sh
#   [ ! -d /srv/vllm2/.venv ] || chown -R $USER:$USER /srv/vllm2/.venv
#   systemctl --user daemon-reload
#   systemctl --user enable --now vllm

[Unit]
Description=vLLM Server -- DiffusionGemma 26B A4B (DGX Spark)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/vllm2

# Overridable via: systemctl --user edit vllm
Environment=PORT=8001
Environment=HOST=0.0.0.0
Environment=MAX_MODEL_LEN=8192
Environment=GPU_MEM_UTIL=0.24
Environment=QUANTIZATION=fp8
Environment=ATTN_BACKEND=TRITON_ATTN
Environment=TOOL_CALL_PARSER=gemma4
Environment=REASONING_PARSER=gemma4
# Server startup script (sources .venv)
ExecStart=/srv/vllm2/vllm-server.sh

# Graceful shutdown -- vLLM needs time to drain
TimeoutStopSec=180
KillSignal=SIGTERM
SendSIGKILL=no

# Auto-restart on crash
Restart=on-failure
RestartSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target