# vllm.service
[Unit]
Description=vLLM Server -- DiffusionGemma 26B A4B (DGX Spark)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/vllm2

# DGX Spark (Grace Blackwell)
# Override using `systemctl --user edit vllm` if needed.
Environment=PORT=8000
Environment=HOST=0.0.0.0
Environment=MAX_MODEL_LEN=262144
Environment=GPU_MEM_UTIL=0.55
Environment=ATTN_BACKEND=TRITON_ATTN
Environment=MODEL_REPO=google/diffusiongemma-26B-A4B-it
Environment=MODEL_ALIAS=google/diffusiongemma-26B-A4B-it
Environment=LLM_ENFORCE_EAGER=false
Environment=MAX_DENOISING_STEPS=48
Environment=CANVAS_LENGTH=256

ExecStart=/srv/vllm2/vllm-server.sh

# vLLM needs drain time on restart/shutdown
TimeoutStopSec=180
KillSignal=SIGTERM
SendSIGKILL=no

Restart=on-failure
RestartSec=20

LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target